# tools/deployment/packaging/kanopya-storage/debian/postinst
#!/bin/bash -e 
echo "ISCSITARGET_ENABLE=true" > /etc/default/iscsitarget
touch /var/log/tftp.log
invoke-rc.d iscsitarget restart
