# tools/deployment/packaging/kanopya-core/debian/postinst
#!/bin/bash -e
update-rc.d kanopya-executor defaults
chmod +x /etc/init.d/kanopya-executor
