# debian/kanopya-monitor.postinst
#!/bin/bash -e
update-rc.d kanopya-grapher defaults
update-rc.d kanopya-collector defaults
chmod +x /etc/init.d/kanopya-grapher
chmod +x /etc/init.d/kanopya-collector
