# debian/kanopya-webui.postinst
#!/bin/bash -e

update-rc.d kanopya-front defaults
chmod +x /etc/init.d/kanopya-front

